# sysutils/trd/pkg/trd.rc
#!/bin/sh
#
# $OpenBSD: trd.rc,v 1.2 2016/12/12 08:02:06 abieber Exp $

daemon="${TRUEPREFIX}/bin/trd"
daemon_flags=""
daemon_user="_tftpd"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_bg=YES

rc_cmd $1

